# systemd/tg-bot.service
[Unit]
Description=Telegram Bot Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/tg-bot-service
EnvironmentFile=/opt/tg-bot-service/.env
ExecStart=/opt/tg-bot-service/venv/bin/python -m bot.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
